# lamp53/conf/memcached/supervisord-memcached.conf
[program:memcached]
command=/usr/bin/memcached -u memcache
numprocs=1
autostart=true
autorestart=true